# templates/smartdns-heartbeat.cron
*/5 * * * * root /usr/local/bin/smartdns-heartbeat
