# templates/smartdns-heartbeat.cron
*/5 * * * * root flock -n /var/run/smartdns-heartbeat.lock -c 'sleep $(awk "BEGIN{srand(); print int(rand()*60)}"); /usr/local/bin/smartdns-heartbeat'

